# systemd/.config/systemd/user/battery-warn.service
[Unit]
Description=Warn on low battery

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/home/anthony/.config/scripts/battery-warn.sh
