# systemd/.config/systemd/user/battery-warn.timer
[Unit]
Description=Warn on low battery

[Timer]
OnCalendar=*-*-* *:0/2:00
AccuracySec=45

[Install]
WantedBy=timers.target
